# zigbee-backup-old-device/scripts/binding-watchdog.service
[Unit]
Description=Zigbee Binding Watchdog - Auto-fixes stale light bindings
Documentation=https://github.com/siva/zigbee
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/zigbee/scripts
ExecStart=/usr/bin/python3 /root/zigbee/scripts/binding-watchdog.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=binding-watchdog

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/root/zigbee/scripts

[Install]
WantedBy=multi-user.target
