# configs/kiosk-toggle/kiosk-toggle.service
[Unit]
Description=Kiosk Toggle Overlay - Floating button to toggle kiosk mode
# Start after display is ready (same pattern as kiosk-browser)
After=default.target display-boot-check.service
Wants=display-boot-check.service

[Service]
Type=simple

# Environment for Wayland
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
# Preload gtk4-layer-shell to fix linking order issue
Environment=LD_PRELOAD=/usr/lib/aarch64-linux-gnu/libgtk4-layer-shell.so.0
# Disable accessibility warnings
Environment=GTK_A11Y=none

# Wait for display to be ready
ExecStartPre=/bin/sleep 8

# Run the GTK4 overlay app
ExecStart=/usr/bin/python3 /opt/kiosk-toggle/kiosk-toggle.py

# Always restart (GTK may exit cleanly even on failure)
# But limit to 5 restarts per 5 minutes to prevent infinite loops
Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
