# configs/pi-reboot/daily-reboot.service
[Unit]
Description=Daily Pi reboot

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -r now
