# backups/configs/wifi-watchdog/wifi-watchdog.service
[Unit]
Description=WiFi Watchdog - Monitor and recover WiFi connection
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-watchdog.sh

# Don't restart on failure - the timer will handle scheduling
Restart=no
