# configs/zigbee-ghost-sweep/zigbee-ghost-sweep-failure.service
[Unit]
Description=Email on zigbee-ghost-sweep failure
# Triggered by OnFailure= in zigbee-ghost-sweep.service.
# Single-shot: reads journalctl + POSTs CRITICAL email via HA API.

[Service]
Type=oneshot
EnvironmentFile=-/opt/zigbee-watchdog/.env
ExecStart=/opt/zigbee-ghost-sweep/notify-failure.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zigbee-ghost-sweep-failure

# Hardening — the notifier only needs: read env, read journal, curl localhost
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
